# sim5360/lib/systemd/system/sim5360.service
[Unit]
Description=sim5360 modem poweron

[Service]
Type=simple
EnvironmentFile=-/etc/default/sim5360
ExecStart=/usr/bin/sim5360 --poweron $OPTIONS
