# sim5360/lib/systemd/system/sim5360.timer
[Unit]
Description=start sim5360 after boot

[Timer]
OnBootSec=10
Unit=sim5360.service

[Install]
WantedBy=multi-user.target
