# adc/DEBIAN/postrm
if [ "$1" = "purge" ] ; then
        systemctl disable adc.service
fi
