# adc/DEBIAN/postinst
systemctl --system daemon-reload
systemctl enable adc.service
systemctl start adc.service
