# sim5360/DEBIAN/postrm
if [ "$1" = "purge" ] ; then
        systemctl disable sim5360.timer
	systemctl disable sim5360_shutdown.service
fi
