# sim5360/DEBIAN/postinst
systemctl enable sim5360.timer
systemctl enable sim5360_shutdown.service
