# powerwatch/DEBIAN/postinst
systemctl --system daemon-reload
systemctl enable powerwatch.service
systemctl start powerwatch.service
