# powerwatch/DEBIAN/postrm
if [ "$1" = "purge" ] ; then
        systemctl disable powerwatch.service
fi
