# one-wire/lib/systemd/system/w1_bus_remove_hook.service
[Unit]
Description=w1 bus remove hook

[Service]
Type=simple
PIDFile=/var/run/w1_bus_remove_hook.pid
ExecStart=/usr/bin/w1_bus_remove_hook -vv -s -p
